# misc/infra/podman/systemd/windmill-phoenix.container
[Container]
Image=ghcr.io/windmill-labs/windmill:main
Name=windmill-phoenix
Pod=phoenix-hydra
# Load secrets from a .env file located in the same directory
EnvironmentFile=.env
# The DATABASE_URL uses the service name 'revenue-db' as the host, not 'localhost'.
# The WINDMILL_DB_PASSWORD will be sourced from the .env file.
Environment="DATABASE_URL=postgresql://postgres:${WINDMILL_DB_PASSWORD}@revenue-db:5432/windmill" "RUST_LOG=info"
Volume=e:/PHOENIXSYSTEM/PhoenixSeed/windmill-scripts:/windmill/scripts:Z
PublishPort=8000:8000
Restart=on-failure
DependsOn=revenue-db.container

[Install]
WantedBy=default.target
