# infra/podman/systemd/phoenix.target
[Unit]
Description=Phoenix Hydra Target
Wants=phoenix-hydra.pod
After=network-online.target

[Install]
WantedBy=default.target
